# agent/weave.service
[Unit]
Description=Weave Agent
Documentation=https://github.com/danohn/weave
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/weave/agent.env
ExecStart=/usr/local/bin/weave
Restart=on-failure
RestartSec=10
# Never give up restarting — network outages can be arbitrarily long
StartLimitIntervalSec=0
# Don't wait 90s for graceful shutdown; our SIGTERM handler is fast
TimeoutStopSec=15

StandardOutput=journal
StandardError=journal
SyslogIdentifier=weave

User=root
# Basic hardening (WireGuard requires root/CAP_NET_ADMIN so full sandboxing isn't possible)
# Note: ProtectHome cannot be used — uv stores the tool venv under /root/.local/
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
