# agent/weave.service
[Unit]
Description=Weave Agent
Documentation=https://github.com/danohn/weave
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/weave/agent.env
ExecStart=/usr/local/bin/weave
Restart=on-failure
RestartSec=10
# Give the network time to settle on boot
ExecStartPre=/bin/sleep 2

StandardOutput=journal
StandardError=journal
SyslogIdentifier=weave

# WireGuard needs root to manage interfaces
User=root

[Install]
WantedBy=multi-user.target
